# root/etc/services.d/duplicacy-web/run
#!/usr/bin/with-contenv bash

cd /app/duplicacy || exit

chmod +x /app/duplicacy/duplicacy-web

exec \
	s6-setuidgid abc ./duplicacy-web
